# unit/custom-conky.timer
[Unit]
Description=Start custom-conky

[Timer]
OnStartupSec= 10

[Install]
WantedBy=timers.target